# deploy/scholar-watch-scrape.service
[Unit]
Description=Scholar Watch Scraper (oneshot)

[Service]
Type=oneshot
User=cornish
Group=cornish
WorkingDirectory=/home/cornish/scholar-watch
EnvironmentFile=/home/cornish/scholar-watch/.env
ExecStart=/home/cornish/scholar-watch/.venv/bin/scholar-watch scrape
TimeoutStartSec=3600
