# deploy/scholar-watch-scrape.timer
[Unit]
Description=Scholar Watch daily scrape timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
